# debian/postinst
#!/bin/bash
set -e

# PostInstall script for taskflow-cli
# This runs after the package is installed

INSTALL_DIR="/usr/share/taskflow-cli"
K8S_DIR="$INSTALL_DIR/k8s"
NAMESPACE="taskflow"

echo "Setting up TaskFlow CLI..."

# Function to check if command exists
command_exists() {
    command -v "$1" >/dev/null 2>&1
}

# Check dependencies
if ! command_exists minikube; then
    echo "ERROR: minikube is not installed. Please install minikube first."
    echo "Visit: https://minikube.sigs.k8s.io/docs/start/"
    exit 1
fi

if ! command_exists kubectl; then
    echo "ERROR: kubectl is not installed. Please install kubectl first."
    echo "Visit: https://kubernetes.io/docs/tasks/tools/"
    exit 1
fi

if ! command_exists docker; then
    echo "ERROR: docker is not installed. Please install docker first."
    echo "Visit: https://docs.docker.com/get-docker/"
    exit 1
fi

# Check if minikube is running
if ! minikube status >/dev/null 2>&1; then
    echo "Starting minikube..."
    minikube start
fi

# Create namespace
echo "Creating namespace: $NAMESPACE"
kubectl create namespace $NAMESPACE --dry-run=client -o yaml | kubectl apply -f -

# Generate secrets.yaml if it doesn't exist
SECRETS_FILE="$K8S_DIR/01-secrets.yaml"
if [ ! -f "$SECRETS_FILE" ]; then
    echo "Generating secrets.yaml..."
    
    # Generate random secrets
    JWT_SECRET=$(openssl rand -base64 32)
    DB_PASSWORD=$(openssl rand -base64 16)
    REDIS_PASSWORD=$(openssl rand -base64 16)
    
    cat > "$SECRETS_FILE" << EOF
apiVersion: v1
kind: Secret
metadata:
  name: taskflow-secrets
  namespace: $NAMESPACE
type: Opaque
stringData:
  jwt-secret: "$JWT_SECRET"
  postgres-password: "$DB_PASSWORD"
  redis-password: "$REDIS_PASSWORD"
EOF
fi

# Pull Docker images with latest tag
echo "Pulling Docker images..."
IMAGES=(
    "your-registry/taskflow-api:latest"
    "your-registry/taskflow-worker:latest"
    "your-registry/taskflow-queue-manager:latest"
)

for image in "${IMAGES[@]}"; do
    echo "Pulling $image..."
    docker pull "$image" || echo "Warning: Failed to pull $image"
done

# Load images into minikube
echo "Loading images into minikube..."
for image in "${IMAGES[@]}"; do
    echo "Loading $image into minikube..."
    minikube image load "$image" || echo "Warning: Failed to load $image"
done

# Apply Kubernetes manifests
echo "Applying Kubernetes manifests..."
kubectl apply -f "$K8S_DIR/infrastructure/" -n $NAMESPACE || true
kubectl apply -f "$K8S_DIR/" -n $NAMESPACE || true

# Wait for pods to be ready
echo "Waiting for pods to be ready..."
kubectl wait --for=condition=ready pod -l app=taskflow-api -n $NAMESPACE --timeout=300s || true

# Setup port forwarding (create systemd service)
echo "Setting up port forwarding service..."

cat > /etc/systemd/system/taskflow-port-forward.service << EOF
[Unit]
Description=TaskFlow Port Forward
After=network.target

[Service]
Type=simple
User=$SUDO_USER
ExecStart=/usr/bin/kubectl port-forward -n $NAMESPACE service/taskflow-api 8080:8080
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
EOF

# Enable and start the service
systemctl daemon-reload
systemctl enable taskflow-port-forward.service
systemctl start taskflow-port-forward.service || echo "Warning: Port forwarding service failed to start"

echo ""
echo "✅ TaskFlow CLI installed successfully!"
echo ""
echo "Quick Start:"
echo "  1. Run 'taskflow' to start the CLI"
echo "  2. Use 'taskflow register' to create an account"
echo "  3. Use 'taskflow login' to authenticate"
echo ""
echo "API is accessible at: http://localhost:8080"
echo "Namespace: $NAMESPACE"
echo ""
echo "To check status: kubectl get pods -n $NAMESPACE"
echo "To view logs: kubectl logs -f deployment/taskflow-api -n $NAMESPACE"
echo ""

exit 0
