# debian/prerm
#!/bin/bash
set -e

# PreRemove script for taskflow-cli
# This runs before the package is removed

echo "Stopping TaskFlow services..."

# Stop port forwarding service
if systemctl is-active --quiet taskflow-port-forward.service; then
    systemctl stop taskflow-port-forward.service
    systemctl disable taskflow-port-forward.service
fi

# Remove systemd service file
rm -f /etc/systemd/system/taskflow-port-forward.service
systemctl daemon-reload

echo "TaskFlow services stopped."

exit 0
